# stigix-all-in-one/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:web-ui]
command=npx tsx server.ts
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/sdwan-traffic-gen/web-ui.err.log
stdout_logfile=/var/log/sdwan-traffic-gen/web-ui.out.log
environment=NODE_ENV="production",PORT="8080"

[program:traffic-gen]
command=/bin/bash /app/traffic-generator.sh client-docker
directory=/app
autostart=true
autorestart=true
startsecs=15
stderr_logfile=/var/log/sdwan-traffic-gen/traffic-gen.err.log
stdout_logfile=/var/log/sdwan-traffic-gen/traffic-gen.out.log

[program:voice-gen]
command=python3 -u voice_orchestrator.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/sdwan-traffic-gen/voice-gen.err.log
stdout_logfile=/var/log/sdwan-traffic-gen/voice-gen.out.log

[program:voice-echo]
command=/bin/bash /app/engines/start_target.sh
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/sdwan-traffic-gen/voice-echo.err.log
stdout_logfile=/var/log/sdwan-traffic-gen/voice-echo.out.log

[program:xfr-target]
command=/usr/local/bin/xfr serve --port %(ENV_XFR_PORT)s --max-duration %(ENV_XFR_MAX_DURATION)ss --rate-limit %(ENV_XFR_RATE_LIMIT)s --allow %(ENV_XFR_ALLOW_CIDR)s
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/sdwan-traffic-gen/xfr-target.err.log
stdout_logfile=/var/log/sdwan-traffic-gen/xfr-target.out.log

[program:mcp-server]
command=python3 -m src.server
directory=/app/mcp-server
autostart=true
autorestart=true
stderr_logfile=/var/log/sdwan-traffic-gen/mcp-server.err.log
stdout_logfile=/var/log/sdwan-traffic-gen/mcp-server.out.log
environment=PYTHONUNBUFFERED="1",MCP_TRANSPORT="sse",MCP_PORT="3100",MCP_HOST="0.0.0.0",CONFIG_DIR="/app/config",LOG_DIR="/var/log/sdwan-traffic-gen",DATA_DIR="/app/mcp-data"
