# stigix-all-in-one/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:web-ui]
command=npx tsx server.ts
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/sdwan-traffic-gen/app.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=2
stdout_logfile=/var/log/sdwan-traffic-gen/app.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=2
environment=NODE_ENV="production",PORT="8080"

[program:traffic-gen]
command=/bin/bash /app/traffic-generator.sh client-docker
directory=/app
autostart=true
autorestart=true
startsecs=15
stderr_logfile=/var/log/sdwan-traffic-gen/app.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=2
stdout_logfile=/var/log/sdwan-traffic-gen/app.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=2

[program:voice-gen]
command=python3 -u voice_orchestrator.py
directory=/app/engines
autostart=true
autorestart=true
stderr_logfile=/var/log/sdwan-traffic-gen/app.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=2
stdout_logfile=/var/log/sdwan-traffic-gen/app.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=2

[program:voice-echo]
command=/bin/bash /app/engines/start_target.sh
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/sdwan-traffic-gen/app.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=2
stdout_logfile=/var/log/sdwan-traffic-gen/app.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=2

[program:xfr-target]
command=/usr/local/bin/xfr serve --port %(ENV_XFR_PORT)s --max-duration %(ENV_XFR_MAX_DURATION)ss --rate-limit %(ENV_XFR_RATE_LIMIT)s --allow %(ENV_XFR_ALLOW_CIDR)s
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/sdwan-traffic-gen/app.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=2
stdout_logfile=/var/log/sdwan-traffic-gen/app.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=2

[program:mcp-server]
command=python3 -m src.server
directory=/app/mcp-server
autostart=true
autorestart=true
stderr_logfile=/var/log/sdwan-traffic-gen/app.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=2
stdout_logfile=/var/log/sdwan-traffic-gen/app.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=2
environment=PYTHONUNBUFFERED="1",MCP_TRANSPORT="sse",MCP_PORT="3100",MCP_HOST="0.0.0.0",CONFIG_DIR="/app/config",LOG_DIR="/var/log/sdwan-traffic-gen",DATA_DIR="/app/mcp-data"

[program:log-forwarder]
command=tail -f /var/log/sdwan-traffic-gen/app.log
autostart=true
autorestart=true
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
